# supervisor/serial_ctrl_bt.conf
[program:serial_ctrl]
command = python -u /opt/serial_tools/control_server.py -d 2
autostart = true
autorestart = true
stopsignal = KILL
