# supervisor/serial_ctrl.conf
[program:serial_ctrl]
command = python -u /opt/serial_tools/control_server.py
autostart = true
autorestart = true
stopsignal = KILL
